# ramdisk-startup.service
[Unit]
Desctiption=Create Ramdisk

[Service]
ExecStart=/usr/bin/sh /srv/createRamdisk.sh

[Install]
WantedBy=multi-user.target